# packages/pinkman/pinkman-0.1.7.tar.gz/pinkman-0.1.7/etc/systemd/system/pinkman_server.service
[Unit]
Description=Pinkman Server
After=multi-user.target
StartLimitIntervalSec=500
StartLimitBurst=5

[Service]
Restart=on-failure
RestartSec=5s
ExecStart=pinkman server run -c /etc/config.pinkman.prod.yml

[Install]
WantedBy=multi-user.target

