# packages/underscorepy/underscorepy-0.6.4.tar.gz/underscorepy-0.6.4/_/supports/systemd/systemd.service
[Unit]
Description={name}
After=network.target

[Service]
Type=simple
ExecStart={cmdline}
User={user}
Restart=on-failure

[Install]
WantedBy=multi-user.target
