# packages/spartaqube/spartaqube-0.1.16.tar.gz/spartaqube-0.1.16/spartaqube/supervisor/supervisor_master.conf
[supervisord]
# environment = HOST_DB_EE=%(ENV_HOST_DB_EE)s, MYSQL_DATABASE=%(ENV_MYSQL_DATABASE)s, MYSQL_USER=%(ENV_MYSQL_USER)s, MYSQL_PASSWORD=%(ENV_MYSQL_PASSWORD)s, SERVICE_DB_EE_SERVICE_HOST=%(ENV_SERVICE_DB_EE_SERVICE_HOST)s
stdout_logfile = /app/supervisor/log/supervisord/supervisord-stdout.log
stderr_logfile = /app/supervisor/log/supervisord/supervisord-stderr.log

[program:spartaqube-daphne]
command = /usr/local/bin/daphne spartaqube.asgi:application -b 0.0.0.0 -p 5000
user = root
directory = /app
autostart = true
autorestart = true
stdout_logfile = /app/supervisor/log/daphne/daphne-stdout.log
stderr_logfile = /app/supervisor/log/daphne/daphne-stderr.log

[program:spartaqube-gunicorn]
# command = /usr/local/bin/gunicorn spartaqube.wsgi:application --bind 0.0.0.0:8000
command = /usr/local/bin/gunicorn spartaqube.wsgi:application -c /app/gunicorn/gunicorn.conf.py
user = root
directory = /app
autostart = true
autorestart = true